# templates/pretalx.socket.j2
[Unit]
Description=pretalx_%I gunicorn socket

[Socket]
ListenStream=/run/gunicorn/pretalx_%i

[Install]
WantedBy=sockets.target
