# templates/pretalx.socket.j2
[Unit]
Description=pretalx{{ pretalx_instance_identifier }} gunicorn socket

[Socket]
ListenStream=/run/gunicorn/pretalx{{ pretalx_instance_identifier }}

[Install]
WantedBy=sockets.target
