# files/pretalx.socket
[Unit]
Description=pretalx gunicorn socket

[Socket]
ListenStream=/run/gunicorn/pretalx

[Install]
WantedBy=sockets.target
